# INSTALL/systemd/exota-milter.service
[Unit]
Description=ExOTA-Milter

[Service]
Restart=always
ExecStart=/usr/loca/sbin/exota-milter.sh

[Install]
WantedBy=multi-user.target
